# scripts/seo-engine/seo-engine.crontab.example
# InstaDeco, moteur SEO sur VPS Hetzner (remplace GitHub Actions seo-engine.yml).
# Installe : voir docs/CRON_VPS_HETZNER.md (section moteur SEO).
# Horaires en UTC (regler le fuseau du VPS si besoin). m h dom mon dow  commande
#
0 6  * * *   /opt/instadeco/run-seo-engine.sh gsc_daily
30 6 * * *   /opt/instadeco/run-seo-engine.sh drift_check
0 7  * * 1   /opt/instadeco/run-seo-engine.sh rank_tracker
0 7  * * 0   /opt/instadeco/run-seo-engine.sh competitor_diff
0 8  1 * *   /opt/instadeco/run-seo-engine.sh citation_batch
