# scripts/seo-engine/seo-engine.crontab.example
# InstaDeco, moteur SEO sur VPS Hetzner (remplace GitHub Actions seo-engine.yml).
# Installe : voir docs/CRON_VPS_HETZNER.md (section moteur SEO).
# Horaires en UTC (regler le fuseau du VPS si besoin). m h dom mon dow  commande
#
# --- Jobs FONCTIONNELS (implementes + testes) ---
0 6  * * *   /opt/instadeco/run-seo-engine.sh gsc_daily      # KPIs Search Console (quotidien)
30 6 * * *   /opt/instadeco/run-seo-engine.sh drift_check    # regressions on-page (quotidien)
0 7  * * 1   /opt/instadeco/run-seo-engine.sh rank_tracker   # positions via GSC (hebdo lundi)
0 8  * * 1   /opt/instadeco/run-seo-engine.sh ctr_optimizer  # auto-optimisation CTR titles (hebdo lundi)
#
# --- Jobs PAS ENCORE PRETS (squelettes) : decommenter quand implementes ---
# 0 7  * * 0   /opt/instadeco/run-seo-engine.sh competitor_diff   # veille concurrents
# 0 8  1 * *   /opt/instadeco/run-seo-engine.sh citation_batch    # GEO, monitoring LLM (payant, phase 3)
